# supervisord.dualstack.conf
[program:dualstack]
directory = /var/www/dualstack/Dualstack-Webbug
environment=PATH="/var/www/dualstack/bin"
user = daviesinc
command=/var/www/dualstack/Dualstack-Webbug/server.py
